# install/etc/profile.d/30-system.sh
export PATH="${HOME}/.local/bin:${HOME}/.local/share/cargo/bin:$PATH"
export EDITOR="nvim"
export VISUAL="bat"
export PAGER="less"
export DIFFPROG="delta"
